# systemd/user/timers.target.wants/screenshotter.timer
[Unit]
Description=Takes a screenshot every half hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
