# systemd/user/appointment.timer
[Unit]
Description=Reminds you of appointments

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
