# systemd/user/prevent_sleep.timer
[Unit]
Description=Prevents the system from going to sleep

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
