# systemd/user/prevent_sleep_youtube.timer
[Unit]
Description=Prevents the system from going to sleep

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
